# etc/udev/rules.d/50-midiconnect.rules
SUBSYSTEM=="usb", DRIVER=="usb" ,ATTR{product}=="USB Audio Device", RUN+="/home/pi/startmidi.sh"
